# rootfs/service/rtorrent/run
#!/bin/sh

exec /usr/local/bin/rtorrent-start
